# server-config/etc/systemd/system/gunicorn-staging.socket
[Unit]
Description=gunicorn-staging socket

[Socket]
ListenStream=/run/gunicorn-staging.sock

[Install]
WantedBy=sockets.target
