# server-config/etc/systemd/system/gunicorn-production.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn-production.sock

[Install]
WantedBy=sockets.target
